# server/bind-stream-deck.service
[Unit]
Description=Bind Elgato Stream Deck using usbip
After=usbipd.service
Requires=usbipd.service

[Service]
ExecStart=/usr/local/bin/bind-stream-deck.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target